# etc/systemd/system/multi-user.target.wants/plex-inotify.service
[Unit]
Description=Plex Inotify daemon

[Service]
Type=exec
ExecStart=/usr/local/bin/plex-inotify.sh

[Install]
WantedBy=multi-user.target
